# install/mimic@.service.in
[Unit]
Description = Start Mimic on %i
After = network.target
Requires = modprobe@mimic.service

[Service]
Type = notify
ExecStart = @@MIMIC_EXEC@@ run %i -F @@MIMIC_CONFIG_PATH@@/%i.conf
Restart = on-failure

User = mimic
RuntimeDirectory = @@MIMIC_RUNTIME_DIR@@
RuntimeDirectoryMode = 0750
RuntimeDirectoryPreserve = yes

CapabilityBoundingSet = CAP_SYS_ADMIN CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW
AmbientCapabilities = CAP_SYS_ADMIN CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW
ProtectSystem = strict

[Install]
WantedBy = multi-user.target
